# files/node_exporter.service
[Unit]
Description=Node Exporter

[Service]
User=prometheus
ExecStart=/Prometheus/node_exporter

[Install]
WantedBy=default.target
